# hass-cloudflare-argo-zerotrust/rootfs/etc/services.d/cloudflared/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Home Assistant Community Add-on: Cloudflare Argo Tunnel
# Take down the S6 supervision tree when Cloudflared fails
# ==============================================================================

if -n { s6-test $# -ne 0 }
if -n { s6-test ${1} -eq 256 }

/run/s6/basedir/bin/halt
