# hass-cloudflare-argo-zerotrust/rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs the Cloudflare daemon
# ==============================================================================

TOKEN=$(bashio::config 'token')

bashio::log.info "Starting cloudflare tunnel."
exec cloudflared tunnel run --token $TOKEN
