# hass-cloudflare-argo/rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs the Cloudflare daemon
# ==============================================================================

bashio::log.info "Starting cloudflare tunnel."
exec cloudflared tunnel run
